# openvpnas_exporter.service
[Unit]
Description=OpenVPN-AS Prometheus Exporter

[Service]
ExecStart=/usr/sbin/openvpnas_exporter

[Install]
WantedBy=multi-user.target